# share/mfemcramfs/tmp/roms/usr/etc/rc.local
#!/bin/bash

/bin/mount -t yaffs /dev/mtdblock3 /usr/yaffs/

